# recipes-connectivity/bcm43xx-utils/bcm43xx-utils/variscite-wifi.service
[Unit]
Description=Variscite WIFI setup service

[Service]
Type=forking
ExecStart=/etc/wifi/variscite-wifi

[Install]
WantedBy=network-online.target
